# appdaemon.service
[Unit]
Description=AppDaemon 4 for Home Assistant
After=network.target

[Service]
Type=simple
User=cciaz
# Point WorkingDirectory at the project directory
WorkingDirectory=/home/cciaz/Desktop/HADock
# Use the startup script that loads environment variables
ExecStart=/home/cciaz/Desktop/HADock/start_appdaemon.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
